# udev-rule-hello-world/usb.rules
ACTION=="add", SUBSYSTEM=="block", ENV{DEVTYPE}=="partition", RUN+="/bin/sh -c 'echo You just plugged a USB stick! >> /home/usb.log'"
